# deploy/nexdownsave.service
[Unit]
Description=NexDownSave Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/venv/bin/python __PROJECT_DIR__/bot.py
ExecStartPre=__PROJECT_DIR__/venv/bin/python __PROJECT_DIR__/healthcheck.py
Restart=always
RestartSec=3
TimeoutStopSec=20
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexdownsave

[Install]
WantedBy=multi-user.target
