# deploy/nexdownsave.service
[Unit]
Description=NexDownSave Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=casperhood
WorkingDirectory=/home/casperhood/.codex/NexDownSave
ExecStartPre=/home/casperhood/.codex/NexDownSave/venv/bin/python /home/casperhood/.codex/NexDownSave/healthcheck.py --allow-missing-db
ExecStart=/home/casperhood/.codex/NexDownSave/venv/bin/python /home/casperhood/.codex/NexDownSave/bot.py
Restart=always
RestartSec=3
TimeoutStopSec=20
Environment=PYTHONUNBUFFERED=1
NoNewPrivileges=yes
PrivateTmp=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexdownsave

[Install]
WantedBy=multi-user.target
